# seedDownloader.ini
[program:seedDownloader]
directory=/root/seedDownloader
command=/root/seedDownloader/seedDownloader
process_name=%(program_name)s-%(process_num)d
username=root
numprocs=4
autostart=true
umask=022
priority=999
redirect_stderr=true
stopwaitsecs = 600
stopasgroup=true
stdout_logfile_maxbytes=10MB
stdout_logfile=/root/seedDownloader/logs/downloader.log